# launchd/com.han.hermes-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.han.hermes-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/han/.local/bin/hermes</string>
    <string>dashboard</string>
    <string>--tui</string>
    <string>--no-open</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>9119</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/han</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/han/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/han</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/han/.hermes/logs/dashboard.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/han/.hermes/logs/dashboard.launchd.err.log</string>
</dict>
</plist>
